# deploy/baremetal/zombie-runner.service
# usezombie Runner — host-resident execution-plane daemon.
# Leases work from the zombied control plane over HTTPS, runs each event in a
# forked sandboxed child, and reports back. Holds zero datastore credentials.
#
# Install:
#   sudo cp zombie-runner.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now zombie-runner

[Unit]
Description=usezombie Runner
Documentation=https://docs.usezombie.com
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/zombie-runner
Restart=always
RestartSec=5
# A stopped runner's in-flight lease is reclaimed by the control plane
# (lease_expires_at + fencing_token), so a bounded stop never loses or
# double-runs work; the timeout just lets an in-flight child finish.
TimeoutStopSec=300
KillMode=mixed

# Required runner config — provisioned per host into the env file below by the
# runner bootstrap playbook (never committed; the runner-token is a secret):
#   ZOMBIE_API_URL       — zombied HTTPS API base, e.g. https://api.usezombie.com
#   ZOMBIE_RUNNER_TOKEN  — pre-minted `zrn_…` token (Option B): the platform
#                          admin mints it via `POST /v1/runners` on the live
#                          control plane and stores it in 1Password under
#                          `…/zombie-…-worker-ant/runner-token`; the daemon
#                          authenticates with it directly and NEVER self-registers
#   RUNNER_HOST_ID       — stable machine identifier (reuse the host's vault
#                          `hostname` field for consistency)
#   (optional) RUNNER_SANDBOX_TIER (prod must override the dev_none default),
#              RUNNER_WORKSPACE_BASE,
#              RUNNER_NETWORK_POLICY (set to registry_allowlist on bare-metal
#                                    hosts that need agent dependency installs;
#                                    default deny_all)
# Bootstrap order (chicken-and-egg): the updated zombied (with the platform-
# admin enrollment gate) must be deployed FIRST; only then can an admin mint a
# `zrn_` for this host. Until that's done, DEV_WORKER_READY stays `false`.
EnvironmentFile=-/etc/default/zombie-runner

# Ensure the sandbox runtime directory exists.
RuntimeDirectory=zombie
RuntimeDirectoryMode=0755

# Security hardening — the runner forks bwrap user-namespace sandboxes per lease
# (no suid needed). Mirrors the isolation the standalone sandbox sidecar carried before
# the engine folded into the runner.
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/run/zombie /sys/fs/cgroup /tmp
ProtectHome=yes

[Install]
WantedBy=multi-user.target
